# metrics/monitor.service
[Unit]
Description=Project Management Dashboard Monitor Service
After=network.target

[Service]
Type=simple
User=vikd
Group=vikd
WorkingDirectory=/Users/Shared/cursor/project_management_dashboard
Environment=PYTHONPATH=/Users/Shared/cursor/project_management_dashboard
ExecStart=/Users/Shared/cursor/project_management_dashboard/.venv/bin/python -m src.monitor
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only

# Resource limits
CPUQuota=25%
MemoryLimit=256M
TasksMax=10

[Install]
WantedBy=multi-user.target
